# deploy/miniapp-api.service
[Unit]
Description=Bit VPN Mini App API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-bot
EnvironmentFile=/opt/vpn-bot/.env
Environment=PATH=/opt/vpn-bot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/vpn-bot/venv/bin/uvicorn api_miniapp:app --host 127.0.0.1 --port 8765
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
