# deploy/support-bot.service
[Unit]
Description=Bit VPN Support Bot (@HelpBit_bot)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-bot
Environment=PATH=/opt/vpn-bot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/vpn-bot/venv/bin/python -m support_bot.run
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
