# deploy/vpn-bot.service
[Unit]
Description=Bit VPN Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-bot
# BOT_TOKEN и остальное читаются из .env через python-dotenv в bot/config/settings.py (файл должен лежать в WorkingDirectory)
# При желании можно продублировать: EnvironmentFile=-/opt/vpn-bot/.env
Environment=PATH=/opt/vpn-bot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/vpn-bot/venv/bin/python run.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
